# templates/default/init.conf.erb
setuid <%= @user %>
setgid <%= @user %>
start on runlevel [2345]
stop on runlevel [016]

respawn

env DIR=/opt/NzbDrone

exec mono $DIR/NzbDrone.exe
